# host/rootfs/image/etc/s6-linux-init/run-image/service/vm-services/template/data/service/xdg-desktop-portal-spectrum-host/run
#!/bin/execlineb -S1
# SPDX-License-Identifier: EUPL-1.2+
# SPDX-FileCopyrightText: 2024-2025 Alyssa Ross <hi@alyssa.is>

importas -i VM VM

export DBUS_SESSION_BUS_ADDRESS unix:path=/run/portal-bus/${VM}

if { mkdir -p /run/vsock/${VM} }
s6-ipcserver-socketbinder -a 0700 /run/vsock/${VM}/vsock_219

# Notify readiness.
if { fdmove 1 3 echo }
fdclose 3

unshare -inu --
s6-envuidgid fs
s6-applyuidgid -Uzu 0
nsenter --preserve-credentials -S0
  --mount=/run/vm/by-id/${VM}/ns/mnt
  --user=/run/vm/by-id/${VM}/ns/user

s6-setuidgid xdp-spectrum-${VM}
setpriv
  --landlock-access fs
  #--landlock-rule path-beneath:read-file,execute:/nix/store
  #--landlock-rule path-beneath:read-file,execute:/usr/bin
  #--landlock-rule path-beneath:read-file,execute:/usr/lib
  #--landlock-rule path-beneath:read-file:/run/vm/by-id/${VM}/portal-bus
  --
xdg-desktop-portal-spectrum-host
