# host/rootfs/usr/lib/systemd/system/cloud-hypervisor@.socket
[Unit]
Description=Cloud Hypervisor control socket for VM %i
PartOf=spectrum-vm@%i.service

[Socket]
ListenStream=/run/vm/by-id/%i/vmm
