# host/rootfs/usr/lib/systemd/system/cloud-hypervisor@.service
[Unit]
Description=VMM for %i
PartOf=spectrum-vm@%i.service

[Service]
Type=simple
ExecStart=/usr/bin/run-vmm
StandardInput=socket
