# host/rootfs/usr/lib/systemd/system/vmm-env.target
[Unit]
After=systemd-sysctl.service kvm.device
Requires=systemd-sysctl.service kvm.device
