# host/rootfs/image/etc/s6-linux-init/run-image/service/vm-services/template/data/service/dbus/run
#!/bin/execlineb -P
# SPDX-License-Identifier: EUPL-1.2+
# SPDX-FileCopyrightText: 2024-2025 Alyssa Ross <hi@alyssa.is>

importas -i VM VM

nsenter --mount=${VM}/mount

unshare
  --cgroup
  --ipc
  --net
  --uts
dbus-daemon
  --config-file /usr/share/dbus-1/session.conf
  --print-address 3
  --address unix:path=${VM}/portal-bus
