# host/rootfs/image/etc/s6-linux-init/run-image/service/vm-services/template/data/service/dbus/run
#!/bin/execlineb -P
# SPDX-License-Identifier: EUPL-1.2+
# SPDX-FileCopyrightText: 2024-2025 Alyssa Ross <hi@alyssa.is>

importas -i VM VM

if {
  redirfd -w 1 data/dbus.conf
  sed "s/@XDP_SPECTRUM_USER@/xdp-spectrum-${VM}/g" /etc/dbus-portal.conf.in
}

s6-ipcserver-socketbinder -B /run/portal-bus/${VM}

fdmove -c 3 0
redirfd -r 0 /dev/null

getcwd -E dir
nsenter --mount=/run/vm/by-id/${VM}/ns/mnt

unshare --cgroup --ipc --net --uts

export LISTEN_FDS 1
getpid LISTEN_PID

dbus-daemon
  --config-file ${dir}/data/dbus.conf
  --print-address 4
  --address systemd:
