# host/rootfs/image/etc/s6-linux-init/run-image/service/s6-linux-init-shutdownd/run
#!/bin/execlineb -WP
# SPDX-License-Identifier: EUPL-1.2+
# SPDX-FileCopyrightText: 2023 Alyssa Ross <hi@alyssa.is>

# cgroups make no sense for shutdownd
s6-linux-init-shutdownd -Bc /etc/s6-linux-init
