# host/rootfs/image/etc/s6-linux-init/run-image/service/vm-services/template/data/service/spectrum-router/run
#!/bin/execlineb -P
# SPDX-License-Identifier: EUPL-1.2+
# SPDX-FileCopyrightText: 2025 Yureka Lilian <yureka@cyberchaos.dev>

importas -i VM VM

bwrap
  --unshare-all
  --unshare-user
  --dev-bind / /
  --setenv RUST_LOG spectrum-router=debug,info
  --tmpfs /tmp
  --dev /dev
  --tmpfs /dev/shm
  --ro-bind /nix /nix
  --ro-bind /etc /etc
  --tmpfs /run
  --ro-bind /usr /usr
  --ro-bind /lib /lib
  --bind /run/vm/by-id/${VM} /run/vm/by-id/${VM}
  --
spectrum-router
  --app-listen-path /run/vm/by-id/${VM}/router-app.sock
  --driver-listen-path /run/vm/by-id/${VM}/router-driver.sock
