# host/rootfs/image/etc/s6-linux-init/run-image/service/vm-services/template/data/service/vhost-user-fs/run
#!/bin/execlineb -WS1
# SPDX-License-Identifier: EUPL-1.2+
# SPDX-FileCopyrightText: 2021-2025 Alyssa Ross <hi@alyssa.is>

s6-ipcserver-socketbinder -a 0700 -B env/virtiofsd.sock

importas -i VM VM
if { chown vmm-${VM} env/virtiofsd.sock }

if { fdmove 1 3 echo }
fdmove -c 3 0
redirfd -r 0 /dev/null

export TMPDIR /run

s6-envuidgid fs
s6-applyuidgid -Uzu 0
importas -i VM VM
nsenter --preserve-credentials -S0
  --mount=/run/vm/by-id/${VM}/ns/mnt
  --user=/run/vm/by-id/${VM}/ns/user

# Show the guest files owned by uid/gid 1000.
unshare -U --map-user 1000 --map-group 1000 --uts --ipc --cgroup

virtiofsd --fd 3 --shared-dir /run/fs/${VM}
